# cron.d/kodi-shutdown
*/4 * * * * root /usr/local/bin/kodi-shutdown-inhibit.sh
